# dvdrip@.service
[Service]
WorkingDirectory=
Type=simple
ExecStart=
User=
